# root/etc/s6-overlay/s6-rc.d/init-branding/up
#!/usr/bin/with-contenv bash

cat <<'EOF'
============================================================
FOSS SmartHome Planner
Bootstrap container for the future planning service
============================================================
EOF
