# root/etc/s6-overlay/s6-rc.d/init-mods-package-install/up
#!/usr/bin/with-contenv bash

echo "Docker Mods hook available. No extra packages requested."
