# root/etc/s6-overlay/s6-rc.d/init-custom-files/up
#!/usr/bin/with-contenv bash

echo "No custom startup files detected."
